# roles/proxmox-common/files/nginx.service.j2
[Unit]
Requires={{ proxmox_proxy_service }}.service
After={{ proxmox_proxy_service }}.service

[Service]
Restart=on-failure
RestartSec=6
