# ansible/roles/proxmox_backup/files/proxmox-backup.timer.j2
[Unit]
Description=Proxmox backup timer

[Timer]
Unit=proxmox-backup.service
OnCalendar={{ backup_schedule }}
WakeSystem=true
Persistent=true

[Install]
WantedBy=timers.target
